# wheelbarrow.service
[Unit]
Description=Wheelbarrow MMO (uvicorn / FastAPI)
After=network-online.target mariadb.service
Wants=network-online.target

[Service]
Type=simple
User=henry
Group=henry
WorkingDirectory=/home/henry/wheelbarrow
EnvironmentFile=/home/henry/wheelbarrow/.env
ExecStart=/home/henry/wheelbarrow/venv/bin/python3 -m uvicorn server.main:app --host 127.0.0.1 --port 8765
Restart=on-failure
RestartSec=5
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target
